# archive/v3.2.0/scripts/docker-cleanup.service
[Unit]
Description=Docker Resource Cleanup Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/home/pete/Documents/Projects/News Intelligence/scripts/docker-cleanup.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target







