# archive/v3.2.0/configs/news-system.service
[Unit]
Description=News Intelligence System v2.0.0
Documentation=https://github.com/your-repo/news-system
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pete/Documents/Projects/News Intelligence
ExecStart=/home/pete/Documents/Projects/News Intelligence/scripts/system-recovery-optimized.sh start
ExecStop=/home/pete/Documents/Projects/News Intelligence/scripts/system-recovery-optimized.sh shutdown
ExecReload=/home/pete/Documents/Projects/News Intelligence/scripts/system-recovery-optimized.sh start
TimeoutStartSec=300
TimeoutStopSec=60

# Restart policy
Restart=on-failure
RestartSec=30

# Security settings
User=pete
Group=pete
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pete/Documents/Projects/News Intelligence

[Install]
WantedBy=multi-user.target
