# news-intelligence-auto-build.service
[Unit]
Description=News Intelligence System Auto Build
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=pete
WorkingDirectory=/home/pete/Documents/Projects/News Intelligence
ExecStart=/home/pete/Documents/Projects/News Intelligence/auto-build-prod.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
