# archive/final-consolidation-20251004/project_cleanup_20250926/old_services/news-intelligence-api.service
[Unit]
Description=News Intelligence API Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=pete
WorkingDirectory=/home/pete/Documents/Projects/News Intelligence/api
Environment=DB_HOST=localhost
Environment=DB_NAME=news_intelligence
Environment=DB_USER=newsapp
Environment=DB_PASSWORD=newsapp_password
Environment=DB_PORT=5432
Environment=REDIS_URL=redis://localhost:6379/0
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
