# infrastructure/newsplatform-backup.cron
# News Intelligence — Database backups on Widow
# Daily at 3 AM, weekly on Sundays at 4 AM
# Paths substituted by setup_widow_app.sh
# (Log archive runs on primary — it has the activity logs)

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

0 3 * * * WIDOW_USER /opt/news-intelligence/scripts/db_backup.sh >> /opt/news-intelligence/logs/backup.log 2>&1
0 4 * * 0 WIDOW_USER /opt/news-intelligence/scripts/db_backup_weekly.sh >> /opt/news-intelligence/logs/backup.log 2>&1
